# supervisord.conf
[supervisord]
nodaemon=true

[program:django]
command=gunicorn MeningitisPredictionProject.wsgi:application --bind 0.0.0.0:${PORT:-8000}
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/django.err.log
stdout_logfile=/var/log/django.out.log

[program:celery_worker]
command=/opt/venv/bin/celery -A MeningitisPredictionProject worker --loglevel=debug
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/celery_worker.err.log
stdout_logfile=/var/log/celery_worker.out.log

[program:celery_beat]
command=/opt/venv/bin/celery -A MeningitisPredictionProject beat --loglevel=debug
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/celery_beat.err.log
stdout_logfile=/var/log/celery_beat.out.log